# supervisord.conf
[supervisord]
nodaemon=true

[program:zookeeper]
directory=/
user=root
command=/kafka_2.12-2.6.0/bin/zookeeper-server-start.sh /kafka_2.12-2.6.0/config/zookeeper.properties
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:kafka]
directory=/
user=root
command=/kafka_2.12-2.6.0/bin/kafka-server-start.sh /kafka_2.12-2.6.0/config/server.properties
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:schema-registry]
directory=/
user=root
command=java -XX:+UnlockExperimentalVMOptions -XX:+UseCGroupMemoryLimitForHeap -XX:+UseG1GC -XX:MaxRAMFraction=2 -jar /schema-registry_6.2.0-0/kafka-schema-registry-package-6.2.0-0-standalone.jar /schema-registry_6.2.0-0/schema-registry.properties
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
